# systemd/.config/systemd/user/talon.service
[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=%h/talon/run.sh

[Install]
WantedBy=graphical-session.target

