# navidrome-status/.config/systemd/user/navidrome-status.service
[Unit]
Description=Update polybar with currently playing music in navidrome

[Service]
ExecStart=%h/bin/navidrome-status
